# krokod/tmux-krokod.service
[Unit]
Description=tmux krokod session with Claude Code
After=network.target
Requires=krotam.service

[Service]
Type=forking
RemainAfterExit=yes
ExecStart=/bin/bash -c 'tmux has-session -t krokod 2>/dev/null || tmux new-session -d -s krokod'
ExecStop=/bin/bash -c 'tmux kill-session -t krokod 2>/dev/null; true'
Environment=HOME=/home/yn441611
Environment=BROWSER=none
Environment=PATH=/home/yn441611/.local/share/fnm/aliases/default/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

[Install]
WantedBy=default.target
